# srcs/phpmyadmin/srcs/supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

[eventlistener:processes]
command=sh -c 'echo "READY"; while read line; do echo "$line" >&2; supervisorctl shutdown; done'
events=PROCESS_STATE_STOPPED, PROCESS_STATE_EXITED, PROCESS_STATE_FATAL

[include]
files = /etc/supervisor.d/*.ini

[program:php-fpm7]
command=/usr/sbin/php-fpm7 -F
autostart=true
autorestart=false

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=false